# services/systemd/argus.service.in
[Unit]
Description=Argus Network Monitor
Documentation=https://github.com/rezdm/Argus.cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/argus -s @CMAKE_INSTALL_PREFIX@/etc/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings
User=argus
Group=argus
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/argus /var/log/argus
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Working directory
WorkingDirectory=@CMAKE_INSTALL_PREFIX@

# Environment
Environment=TMPDIR=/var/lib/argus/tmp

[Install]
WantedBy=multi-user.target